# debian/crowdsec-firewall-bouncer.postrm
systemctl stop crowdsec-firewall-bouncer
systemctl disable crowdsec-firewall-bouncer
