# debian/crowdsec-firewall-bouncer-iptables.prerm
#!/bin/sh

set -eu

BOUNCER="crowdsec-firewall-bouncer"

systemctl stop "$BOUNCER" || echo "cannot stop service"
systemctl disable "$BOUNCER" || echo "cannot disable service"
